# usr/lib/systemd/system/auditd.service.d/restarts.conf
[Unit]
RefuseManualStop=no

[Service]
ExecStartPost=-/usr/bin/auditctl -R /etc/audit/audit.rules
ExecStopPost=/usr/bin/auditctl -R /etc/audit/audit-stop.rules
