# usr/lib/systemd/system/headerdev.service
[Unit]
Description=Decrypt disks
ConditionPathExists=/etc/initrd-release
RequiresMountsFor=/headerdev
DefaultDependencies=no
Before=cryptsetup-pre.target

[Service]
Type=oneshot
ExecStart=/usr/share/ppz/bin/headerdev

[Install]
RequireddBy=cryptsetup.target
