# usr/lib/systemd/system/keep-net-up.service
[Unit]
Description=Keep-Net-Up

[Service]
ExecStart=/usr/bin/keep-net-up
