# etc/systemd/system/keep-net-up.service.d/example.conf
# you probably want to copy this to local.conf or something
# before using it. it's not in pacman's backup array...

[Service]
#ExecStart=
